# contrib/systemd/awob.service
# awob — systemd user service unit.
#
# Install:
#   $ mkdir -p ~/.config/systemd/user
#   $ cp contrib/systemd/awob.service ~/.config/systemd/user/
#   $ systemctl --user daemon-reload
#   $ systemctl --user enable --now awob.service
#
# Or, if your compositor uses uwsm / systemd graphical-session targets
# (e.g. Hyprland with uwsm, gnome-session, etc.), the
# `WantedBy=graphical-session.target` already wires it into your
# session lifecycle. The unit will start when you log in and stop
# cleanly when you log out.
#
# Override the theme, themes dir, or socket via a drop-in:
#   $ systemctl --user edit awob.service
#   [Service]
#   Environment=AWOB_THEME=tinct
#
# This unit deliberately doesn't restart on *clean* exit (Type=simple
# + Restart=on-failure) — the daemon owns the IPC socket, and a quick
# respawn loop on missing-Wayland would just spam the journal.

[Unit]
Description=Another Wayland Overlay Bar
Documentation=https://github.com/jmylchreest/awob
PartOf=graphical-session.target
After=graphical-session.target
Requires=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=simple
# Defensive cleanup — if a previous instance was SIGKILL'd or the
# system was hard-rebooted, the socket file lingers and the next
# bind() fails with EADDRINUSE. The leading `-` makes the failure
# tolerated (no socket present is fine). The daemon also unlinks
# stale sockets internally, but the in-process cleanup needs write
# access to the runtime dir (granted by the broader ReadWritePaths
# below) — this ExecStartPre is the belt to that braces.
ExecStartPre=-/usr/bin/rm -f %t/awob.sock
ExecStart=%h/.cargo/bin/awob-daemon
Restart=on-failure
RestartSec=2

# Logging goes to the user journal. Use `journalctl --user -u awob`.
StandardOutput=journal
StandardError=journal

# Quality-of-life sandboxing. None of these block awob's normal
# operation but reduce blast radius if the daemon misbehaves.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
# Grant write access to the runtime dir, not just the socket file.
# Unlinking + rebinding the socket needs write+exec on the parent
# directory, which `%t/awob.sock` (file-only) doesn't grant.
ReadWritePaths=%t/
# The backlight listener reads /sys/class/backlight, so we can't hide
# /dev. systemd unit files don't support inline comments — keep the
# rationale on its own line above the directive.
PrivateDevices=no
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

[Install]
WantedBy=graphical-session.target
